# systemd/vm-hid.service
[Unit]
Description=VM-HID server

[Service]
ExecStart=/usr/local/bin/vm-hid

[Install]
WantedBy=multi-user.target
