# systemd/vm-hid.socket
[Unit]
Description=VM-HID server socket

[Socket]
ListenStream=/run/vm-hid.sock

[Install]
WantedBy=default.target
